# database/crontab
# run every day at midnight
0 22 * * * postgres /docker-entrypoint-initdb.d/update/update_db.sh -o
# An empty line is required at the end of this file for a valid cron file.

